# packages/DEBIAN/prerm
#!/bin/sh
set -e

echo "Pre-Removal Macro"

if [ "$1" = "remove" ]; then
	systemctl disable localdev || true
	systemctl stop localdev || true

	rm -f __PACKAGE_DIR__/config/settings.cfg
fi

exit 0
